# example/journalcheck.timer
[Unit]
Description=Hourly check of Journal entries

[Timer]
OnCalendar=hourly
