# example/journalcheck.service
[Unit]
Description=Journal Check

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/journalcheck
ExecStart=/opt/journalcheck/bin/journalcheck
